# projectmanager.service
[Unit]
Description=Project Manager Web Server
After=network.target
After=postgresql.service

[Service]
User=tabula
WorkingDirectory=/home/tabula/go/src/github.com/denisbakhtin/projectmanager
ExecStart=/home/tabula/go/src/github.com/denisbakhtin/medical/projectmanager-go -mode=release
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target

